# systemd/agent-builder.service
[Unit]
Description=OODATCAA Agent Daemon - Builder
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/storage/Code/MCPLocalLLM
ExecStart=/usr/bin/python3 %h/storage/Code/MCPLocalLLM/scripts/agent-daemon.py --role builder --interval 60 --owner agent-daemon-builder
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=default.target

